# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Contacts and Calendar synchronization service

[Service]
Type=oneshot
ExecStart=/home/igneo676/.bin/vdirsyncer-guard
