# .config/aconfmgr/files/etc/systemd/system/home-igneo676-Music.automount
[Unit]
Description=Automount for Music SSHFS

[Automount]
Where=/home/igneo676/Music
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
